# swarming_bot/files/swarming_bot.conf
{% from "swarming/map.jinja" import swarming with context %}

description "Swarming bot"

start on filesystem or runlevel [2345]
stop on shutdown

respawn
respawn limit 0 10

setuid {{ swarming.user }}

exec /usr/bin/env python {{ swarming.bot_path }} start_bot
